# capture_and_upload.cron
0,15,30,45 * * * * /home/pi/crackcam_uploader/capture_and_upload.sh >> /home/pi/crackcam_uploader/capture_and_upload.log